# crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0,30 * * * * root /usr/local/bin/dynupdate 
0 1 * * * root /usr/local/bin/certbot
